# agents/spousal-pr-agent/com.spousal-pr.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spousal-pr.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zeeshanhabib/cursor/agents/spousal-pr-agent/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zeeshanhabib/cursor/agents/spousal-pr-agent/check-ins/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zeeshanhabib/cursor/agents/spousal-pr-agent/check-ins/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>

<!--
INSTALL (one-time):

  cp /Users/zeeshanhabib/cursor/agents/spousal-pr-agent/com.spousal-pr.weekly.plist ~/Library/LaunchAgents/
  launchctl load ~/Library/LaunchAgents/com.spousal-pr.weekly.plist

Test manually:

  bash /Users/zeeshanhabib/cursor/agents/spousal-pr-agent/run.sh

Uninstall:

  launchctl unload ~/Library/LaunchAgents/com.spousal-pr.weekly.plist
  rm ~/Library/LaunchAgents/com.spousal-pr.weekly.plist
-->
