# agents/spousal-pr-agent/com.spousal-pr.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spousal-pr.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zeeshanhabib/cursor/agents/spousal-pr-agent/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zeeshanhabib/cursor/agents/spousal-pr-agent/check-ins/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zeeshanhabib/cursor/agents/spousal-pr-agent/check-ins/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>

<!--
INSTALL (one-time) — run from repo root (cursor/):

  mkdir -p ~/Library/LaunchAgents
  cp agents/spousal-pr-agent/com.spousal-pr.weekly.plist ~/Library/LaunchAgents/
  launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.spousal-pr.weekly.plist

Verify:

  launchctl list | grep spousal

Test manually:

  cd /Users/zeeshanhabib/cursor
  FORCE_RUN=1 bash agents/spousal-pr-agent/run.sh

Reinstall (after editing plist):

  launchctl bootout gui/$(id -u)/com.spousal-pr.weekly
  cp agents/spousal-pr-agent/com.spousal-pr.weekly.plist ~/Library/LaunchAgents/
  launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.spousal-pr.weekly.plist

Uninstall:

  launchctl bootout gui/$(id -u)/com.spousal-pr.weekly
  rm ~/Library/LaunchAgents/com.spousal-pr.weekly.plist

TROUBLESHOOTING:
  - cp: No such file or directory → mkdir -p ~/Library/LaunchAgents first
  - Bootstrap failed: 5 → already loaded; check launchctl list

Full docs: agents/CLAUDE.md
-->
