# agents/portfolio-agent/com.portfolio.surveillance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.portfolio.surveillance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zeeshanhabib/cursor/agents/portfolio-agent/run.sh</string>
    </array>

    <!-- Runs Mon-Fri at 8:30am local time (Toronto ET) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/zeeshanhabib/cursor/agents/portfolio-agent/briefings/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zeeshanhabib/cursor/agents/portfolio-agent/briefings/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>

<!--
INSTALL (one-time) — run from repo root (cursor/):

  mkdir -p ~/Library/LaunchAgents
  cp agents/portfolio-agent/com.portfolio.surveillance.plist ~/Library/LaunchAgents/
  launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.portfolio.surveillance.plist

Verify:

  launchctl list | grep portfolio

Test manually:

  cd /Users/zeeshanhabib/cursor
  FORCE_RUN=1 bash agents/portfolio-agent/run.sh

Reinstall (after editing plist):

  launchctl bootout gui/$(id -u)/com.portfolio.surveillance
  cp agents/portfolio-agent/com.portfolio.surveillance.plist ~/Library/LaunchAgents/
  launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.portfolio.surveillance.plist

Uninstall:

  launchctl bootout gui/$(id -u)/com.portfolio.surveillance
  rm ~/Library/LaunchAgents/com.portfolio.surveillance.plist

TROUBLESHOOTING:
  - cp: No such file or directory → mkdir -p ~/Library/LaunchAgents first
  - Bootstrap failed: 5 → already loaded; check launchctl list

NOTE: Mac must be awake at 8:30am. If asleep, launchd runs when it wakes.
Full docs: agents/CLAUDE.md
-->
