# agents/portfolio-agent/com.portfolio.surveillance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.portfolio.surveillance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/administrator/agents/portfolio-agent/run.sh</string>
    </array>

    <!-- Runs Mon-Fri at 8:30am local time (Toronto ET) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/administrator/agents/portfolio-agent/briefings/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/agents/portfolio-agent/briefings/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>

<!--
INSTALL INSTRUCTIONS (one-time setup):
Run these two commands in your terminal:

  cp ~/agents/portfolio-agent/com.portfolio.surveillance.plist ~/Library/LaunchAgents/
  launchctl load ~/Library/LaunchAgents/com.portfolio.surveillance.plist

To test it manually right now (without waiting for 8:30am):
  bash ~/agents/portfolio-agent/run.sh

To uninstall:
  launchctl unload ~/Library/LaunchAgents/com.portfolio.surveillance.plist
  rm ~/Library/LaunchAgents/com.portfolio.surveillance.plist

NOTE: Your Mac must be awake at 8:30am for launchd to fire.
If it's asleep, the job runs as soon as it wakes up that day.
-->
